# templates/etc_containers_systemd_castle_pull.container.j2
[Unit]
Description=Pull castle configuration

[Container]
Image=docker.io/adborden/castle:latest

[Service]
Type=oneshot
